# build/ansible/roles/redis/templates/upstart_redis.conf.j2
description "redis server"
start on runlevel [2345]
stop on shutdown
respawn
respawn limit 30 3
exec sudo -u redis /usr/local/bin/redis-server /etc/redis/redis.conf

